# data/laravel-worker.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work sqs --queue=RoyaltyReport --tries=3 --daemon
autostart=true
autorestart=true
# user=ubuntu
numprocs=2
redirect_stderr=true
stdout_logfile=/app/queue.log​
